# docker/ops/packer/upstart/docker_app.conf
description "Docker app"

start on vagrant-ready or runlevel [2345]
stop on runlevel [!2345]

respawn

console log

script
    sudo service apache2 restart
end script
